# contribs/systemd/nats-jwt-callout.service
[Unit]
Description=NATS auth callout service (OIDC token authentication)
Documentation=https://github.com/sylr/nats-jwt-callout
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nats-jwt-callout
Group=nats-jwt-callout
ExecStart=/usr/bin/nats-jwt-callout -config /etc/nats-jwt-callout/config.yaml
# `systemctl reload nats-jwt-callout` hot-reloads the policy and issuers (SIGHUP).
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
# Stateless network service: ProtectSystem=strict already makes the whole
# filesystem (incl. /etc, which holds the secret-bearing config) read-only, and
# the service writes nothing to disk — so no ReadWritePaths/StateDirectory.

[Install]
WantedBy=multi-user.target
